# roles/webhook/templates/webhook.service.j2
[Unit]
Description=Webhooks

[Service]
ExecStart=/home/ubuntu/PaaSeasy/webhook -hooks  /home/ubuntu/goapp/PaaSeasy/configfiles/hooks.json -hotreload -verbose

[Install]
WantedBy=multi-user.target
