# roles/goapplication/templates/goapplication.service.j2
[Unit]
Description=Webhooks

[Service]
ExecStart=/home/ubuntu/PaaSeasy/paaseasy

[Install]
WantedBy=multi-user.target

