# roles/goapplication/templates/goapplication.service.j2
[Unit]
Description=Go PaaS Application

[Service]
WorkingDirectory=/home/ubuntu/goapp/PaaSeasy
ExecStart=/home/ubuntu/goapp/PaaSeasy/paaseasy

[Install]
WantedBy=multi-user.target

